# config/crontab
00 6 * * * /bin/bash /root/rba-timelapse/1hrrepeat.sh >> /root/process.log